# etc/supervisord.d/lb-update.ini
[program:lb-update]
command=consul-template -consul localhost:8500 -wait 5s -template "/etc/consul-templates/lb.ctmpl:/etc/nginx/conf.d/es-lb.conf:nginx -s reload"
redirect_stderr=true
autostart=true
startsecs=10
stdout_logfile=/var/log/supervisor/%(program_name)s.log
